# launchd/com.hermes-sidecar.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  hermes-sidecar LaunchAgent for macOS.

  Install:  cp launchd/com.hermes-sidecar.plist ~/Library/LaunchAgents/
            launchctl load ~/Library/LaunchAgents/com.hermes-sidecar.plist

  Uninstall: launchctl unload ~/Library/LaunchAgents/com.hermes-sidecar.plist
             rm ~/Library/LaunchAgents/com.hermes-sidecar.plist

  Logs:     ~/.hermes/sidecar/daemon.log
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hermes-sidecar</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/hermes-sidecar</string>
        <string>daemon</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/REPLACE_WITH_USERNAME/.hermes/sidecar/stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/REPLACE_WITH_USERNAME/.hermes/sidecar/stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/REPLACE_WITH_USERNAME</string>
    </dict>

    <key>ProcessType</key>
    <string>Background</string>

    <key>Nice</key>
    <integer>10</integer>
</dict>
</plist>
